# chroot/etc/xdg/autostart/cdrom.desktop
[Desktop Entry]
Exec=sh -c "gio mount -d /dev/sr0 && xdg-open /media/*/*"
