# chroot/etc/xdg/autostart/cdrom.desktop
[Desktop Entry]
Name=Automatically show CD-ROM contents
Exec=sh -c "gio mount -d /dev/sr0 && xdg-open /media/*/* || :"
Type=Application
